# system_env/plant_monitor-logger.service
# /etc/systemd/system/plant_monitor-logger.service

[Unit]
Description=Plant Monitor Data Logger Daemon
After=network.target

[Service]
# このサービスを実行するユーザー (piユーザーを想定)
User=pi
# グループ
Group=pi

# 実行するコマンドのフルパス
# /home/pi/plant_dashboard/.venv/bin/python はご自身の環境に合わせて変更してください
ExecStart=/home/pi/plant_dashboard/.venv/bin/python /home/pi/plant_dashboard/plant_monitor_daemon.py

# 作業ディレクトリ
WorkingDirectory=/home/pi/plant_dashboard

# エラー発生時に自動で再起動する
Restart=always
# 再起動までの待機時間
RestartSec=10

# 標準出力と標準エラー出力をシステムログ(journald)に送る
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target